# gmail-lead-sync.service
[Unit]
Description=Gmail Lead Sync & Response Engine
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=gmail-sync
Group=gmail-sync
WorkingDirectory=/opt/gmail-lead-sync
ExecStart=/opt/gmail-lead-sync/venv/bin/python -m gmail_lead_sync start
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment variables
Environment="ENCRYPTION_KEY=REPLACE_WITH_YOUR_ENCRYPTION_KEY"
Environment="DATABASE_URL=sqlite:////var/lib/gmail-lead-sync/gmail_lead_sync.db"
Environment="LOG_LEVEL=INFO"

# Optional: Use environment file instead of inline variables
# EnvironmentFile=/etc/gmail-lead-sync/environment

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/gmail-lead-sync /var/log/gmail-lead-sync
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictAddressFamilies=AF_INET AF_INET6

[Install]
WantedBy=multi-user.target
